# beaglebone-master-config-files/etc/supervisor/conf.d/master_server.conf
[program:master_server]
command=python2 master_server.py
num_procs=1
directory=/home/debian/dreamscape2/Dreamscape-II/DreamscapeII/roles/server/files/
user=root
autorestart=true
redirect_stderr=true
