# beaglebone-master-config-files/etc/supervisor/conf.d/picontrol.conf
[program:picontrol]
command=node app
num_procs=1
directory=/opt/piControl/
user=root
autorestart=true
redirect_stderr=true
